# sv/artix-oneshot-rc/finish
#!/bin/sh
. /usr/lib/rc/functions
stop_oneshot_services
