# sv/artix-oneshot-rc/run
#!/bin/sh
. /usr/lib/rc/functions
run_oneshot_services
exec chpst -b artix-oneshot pause
