# sv/udevd/finish
#!/bin/sh
[ -S /run/udev/control ] && rm -f /run/udev/control
